# archiso-profile/airootfs/etc/systemd/system/erinos-update.service
[Unit]
Description=ErinOS daily system update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/erinos update
User=root
