# archiso-profile/airootfs/etc/systemd/system/erinos-update.timer
[Unit]
Description=ErinOS daily update timer

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
